# deploy/supervisord.conf
[supervisord]
nodaemon=false
logfile=%(here)s/../logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(here)s/../supervisord.pid

[unix_http_server]
file=%(here)s/../supervisord.sock

[supervisorctl]
serverurl=unix://%(here)s/../supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:openmoss]
command=%(here)s/../.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 6565
directory=%(here)s/..
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=%(here)s/../logs/stdout.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/stderr.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1"
